# supervisord.conf
[supervisord]
nodaemon=true

[program:docker-gen]
command=docker-gen -watch -notify "kill -HUP $(pidof tor)" /app/torrc.tmpl /etc/tor/torrc

[program:tor]
command=tor -f /etc/tor/torrc